# pure-ftpd/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=false

[program:pureftpd]
command=pure-ftpd -lpuredb:/etc/pureftpd.pdb
environment=STNORESTART=yes
autorestart=true
redirect_stderr=true


