# Syncthing/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=false

[program:syncthing]
command=/usr/local/bin/syncthing -home=/mnt/media/.syncthing
environment=STNORESTART=yes
autorestart=true
redirect_stderr=true
